# labs/Arduino/docker-helper/udev-rules/99-docker-usbdev.rules
# for USB devices

TAG!="uaccess", GOTO="docker-usbdev_end"

ACTION=="remove", SUBSYSTEM=="usb", RUN+="/opt/bin/docker-usbdev.sh"
ACTION=="add", SUBSYSTEM=="usb", RUN+="/opt/bin/docker-usbdev.sh"

LABEL="docker-usbdev_end"
